# systemd/csl-export.service
[Unit]
Description=CSL Dray Export Monitor (JsonCargo API tracking)
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/root/csl-bot
ExecStart=/usr/bin/python3 export_monitor.py
EnvironmentFile=/root/csl-bot/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csl-export

[Install]
WantedBy=multi-user.target
