# systemd/csl-health-check.service
[Unit]
Description=CSL Health Check
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/root/csl-bot
ExecStart=/usr/bin/python3 health_check.py
EnvironmentFile=/root/csl-bot/.env
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csl-health-check
