# systemd/csl-gmail-monitor.service
[Unit]
Description=CSL Gmail Monitor (5-min polling)
After=network.target postgresql.service

[Service]
Type=simple
WorkingDirectory=/root/csl-bot/csl-doc-tracker
ExecStart=/usr/bin/python3 gmail_monitor.py
EnvironmentFile=/root/csl-bot/csl-doc-tracker/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csl-gmail-monitor

[Install]
WantedBy=multi-user.target
