# systemd/csl-webhook.service
[Unit]
Description=CSL Macropoint Webhook Receiver (port 5000)
After=network.target

[Service]
Type=simple
WorkingDirectory=/root/csl-bot
ExecStart=/usr/bin/python3 webhook.py
EnvironmentFile=/root/csl-bot/.env
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=csl-webhook

[Install]
WantedBy=multi-user.target
